# prod-deployment-scripts/netly-app.service
[Unit]
Description=Netly Application

After=network.target nginx.service postgresql.service
Requires=nginx.service postgresql.service

[Service]
Restart=always
RestartSec=5
User=root
WorkingDirectory=/opt/app/netly

# JVM memory flags for netly (2GB heap, conservative allocation)
ExecStart=/usr/bin/java \
  -Xmx2g \
  -Xms1g \
  -XX:MaxDirectMemorySize=512m \
  -XX:MaxMetaspaceSize=256m \
  -XX:MetaspaceSize=128m \
  -XX:+UseG1GC \
  -XX:MaxGCPauseMillis=300 \
  -XX:G1HeapRegionSize=8m \
  -Xss256k \
  -Xlog:gc*:file=/var/log/netly/gc.log:time,uptime:filecount=5,filesize=10m \
  -jar /opt/app/netly/netly.jar

SuccessExitStatus=143

# Systemd memory limits
MemoryLimit=3G
MemoryHigh=2.5G
CPUQuota=200%

# OOM protection
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target
